# systemd/pm-monitor.service
[Unit]
Description=pm-system monitor
After=network.target

[Service]
Type=simple
User=pm
WorkingDirectory=/home/pm/pm-system
EnvironmentFile=/home/pm/pm-system/.env
ExecStart=/home/pm/.venv/bin/python monitor.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
